# packages/aiy-voice-services/debian/aiy_voice_classic.service
[Unit]
Description=AIY Voice HAT overlay

[Service]
Type=oneshot
ExecStart=-/opt/aiy/voice-services/apply_voice_classic_overlay.sh
RemainAfterExit=no
StandardOutput=journal

[Install]
WantedBy=multi-user.target
